# config/includes.chroot/etc/skel/.config/autostart/bg.desktop
[Desktop Entry]
Type=Application
Name=BGSET
Exec=sh -c 'gsettings set org.mate.background picture-filename /etc/bg/brew.png'